# 3rd_party/supervisor/device-server/etc/supervisor/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[inet_http_server]          ; inet (TCP) server disabled by default
port=127.0.0.1:9001         ; ip_address:port specifier, *:port for all iface
;username=user              ; default is no username (open server)
;password=123               ; default is no password (open server)

[supervisord]
logfile=/opt/latency_server/log/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:latency_server]
command=ti-device-server --host '0.0.0.0' --port 15003 --warmup 50 --repeat 50 --number 25 --working-dir /opt/latency_server/working_dir --reboot-after-measure
process_name=%(program_name)s ; process_name expr (default %(program_name)s)
numprocs=1                    ; number of processes copies to start (def 1)
autostart=true                ; start at supervisord start (default: true)
startsecs=1                   ; # of secs prog must stay up to be running (def. 1)
startretries=3                ; max # of serial start failures when starting (default 3)
autorestart=unexpected        ; when to restart if exited after running (def: unexpected)
exitcodes=0                   ; 'expected' exit codes used with autorestart (default 0)
stopsignal=TERM               ; signal used to kill process (default TERM)
stopwaitsecs=30               ; max num secs to wait b4 SIGKILL (default 10)
redirect_stderr=true          ; redirect proc stderr to stdout (default false)
stdout_logfile=/opt/latency_server/log/latency_server.log        ; stdout log path, NONE for none; default AUTO
stdout_logfile_maxbytes=10MB  ; max # logfile bytes b4 rotation (default 50MB)
stdout_logfile_backups=10     ; # of stdout logfile backups (0 means none, default 10)
stdout_syslog=false           ; send stdout to syslog with process name (default false)
environment=TIDL_TOOLS_PATH="/opt/latency_server/tidl_tools",TIDL_RT_PERFSTATS="1",PYTHONUNBUFFERED="1"        ; process environment additions (def no adds)
